# crontab
5 16 * * * source /home/ec2-user/Python-3.11.4/env/bin/activate && python /home/ec2-user/Grahams-Emails/lambda_function.py >> /home/ec2-user/Grahams-Emails/output.txt
